# etc/cron/backup_vms.cron
#Min    Hrs     DoM     Mon     DoW     User    Script
# If it's the first saturday of the month, run everything!!!!!
0       2       *       *       6       root    [ $(date +\%d) -le 07 ] && /root/ovirt-backup/etc/backup/monthly && /root/ovirt-backup/etc/backup/weekly && /root/ovirt-backup/etc/backup/daily

# If it's any other saturday of the month, just run weekly and daily
0       2       *       *       6       root    [ $(date +\%d) -ge 08 ] && /root/ovirt-backup/etc/backup/weekly && /root/ovirt-backup/etc/backup/daily 

# For all days other than saturday, just run the daily backup
0       2       *       *       0,1,2,3,4,5     root    /root/ovirt-backup/etc/backup/daily
